# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
#sleep 3 && /etc/init.d/radiusd restart
#/etc/init.d/chilli restart
sleep 10 && sync_time.sh quiz.int.vidio.com
exit 0
